# astrod_sanctuary_kit/forgejo.service
[Unit]
Description=Forgejo (council Git forge)
After=network.target

[Service]
Type=simple
User=council
WorkingDirectory=/srv/forgejo
ExecStart=/srv/forgejo/forgejo web
Restart=on-failure
Environment=USER=council HOME=/home/council

[Install]
WantedBy=multi-user.target
